# ansible/roles/bgp/templates/bsp-bgp-check.timer.j2
[Unit]
Description=BGP health-check timer for shard-proxy

[Timer]
OnBootSec=30s
OnUnitActiveSec=10s
Unit=bsp-bgp-check.service

[Install]
WantedBy=timers.target
